# rootfs/etc/cron.weekly/zfs-snap
#!/bin/sh

set -e

date=$(date -u +%Y-%m-%d-%H%M)

ZFS_ROOT=$(awk '$2 == "/" && $3 == "zfs" {print $1}' /proc/mounts)
ZFS_HOME=$(awk '$2 == "/home" && $3 == "zfs" {print $1}' /proc/mounts)

[ -n "$ZFS_ROOT" ]
[ -n "$ZFS_HOME" ]

zfs snapshot "$ZFS_ROOT"@auto_weekly_"$date"
zfs snapshot "$ZFS_HOME"@auto_weekly_"$date"

{
  zfs list -t snapshot -s creation -o name -H "$ZFS_ROOT" | grep @auto_weekly_ | head -n -4
  zfs list -t snapshot -s creation -o name -H "$ZFS_HOME" | grep @auto_weekly_ | head -n -4
} |
  while IFS= read -r snap; do
    zfs destroy "$snap" || continue
  done
