# rootfs/etc/cron.monthly/fstrim
#!/bin/sh
fstrim -v /boot/efi | logger -t fstrim -p daemon.info
